# config/crontab/test
1 * * * * echo test >> /var/log/cron.log 2>&1
